# supervisor.conf
[program:website]
command=---SITE_PATH---/venv/bin/python tornado_server.py
directory=---SITE_PATH---
user=---WEB_USER---
