# supervisor.conf
[program:website]
command=venv/bin/python tornado_server.py
directory=---SITE_PATH---
user=---WEB_USER---
